# VM_files/start_CTF.service
[Unit]
After=docker.service

[Service]
ExecStart=/home/kaguya/ProjectKaguya/VM_files/start_CTF.sh

[Install]
WantedBy=default.target
